# systemd/system-backup.service
[Unit]
Description=Daily System Backup


[Service]
Type=oneshot
ExecStart=/usr/local/bin/system-backup.sh
