# systemd/system-backup.timer
[Unit]
Description=Run System Backup Daily


[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true


[Install]
WantedBy=timers.target




